# ansible/roles/web-server/templates/iodine.service.j2
[Unit]
Description=Iodine HTTP Server
After=network.target

[Service]
User=deploy
UMask=0000
EnvironmentFile=-/etc/environment

WorkingDirectory=/var/www/ftc_results/current/
ExecStart=/usr/local/bin/bundle exec iodine -p 0 -b /var/www/ftc_results/shared/tmp/sockets/iodine.sock

TimeoutSec=15
Restart=always
KillMode=process

[Install]
WantedBy=multi-user.target